# packages/workbooksd/release/installer/scripts/postinstall
#!/bin/sh
# Postinstall — runs as ROOT after the .pkg lays down:
#
#   /Applications/Workbooks.app   (Tauri Manager UI, opened on demand)
#   /usr/local/bin/workbooksd     (the daemon binary, supervised below)
#
# Four responsibilities:
#
#   1. Bootstrap the user-scoped LaunchAgent so the daemon runs at
#      login + restarts on crash. The user never has to open
#      Workbooks.app for the daemon to start.
#
#   2. Wipe stale runtime.json from any previous (possibly crashed)
#      daemon so the freshly-bootstrapped one writes its own port
#      cleanly.
#
#   3. Refresh Launch Services so /Applications/Workbooks.app's UTI
#      declarations land without a logout.
#
#   4. Set Workbooks as the default handler for sh.workbooks.workbook
#      so DOUBLE-clicking a .workbook.html (not just right-click)
#      routes to Workbooks. Without this, macOS only goes to a
#      UTI-specific app when the user has explicitly set the per-UTI
#      default; otherwise it falls back to the broader public.html
#      handler (Safari).
#
# Idempotent — safe to re-run on upgrade.

set -e

CONSOLE_UID=$(stat -f%u /dev/console)
CONSOLE_USER=$(id -un "$CONSOLE_UID" 2>/dev/null || echo root)
CONSOLE_HOME=$(dscl . -read /Users/"$CONSOLE_USER" NFSHomeDirectory 2>/dev/null | awk '{print $2}')
[ -z "$CONSOLE_HOME" ] && CONSOLE_HOME="/Users/$CONSOLE_USER"

BIN="/usr/local/bin/workbooksd"
PLIST_DIR="$CONSOLE_HOME/Library/LaunchAgents"
PLIST="$PLIST_DIR/sh.workbooks.workbooksd.plist"
LOG_DIR="$CONSOLE_HOME/Library/Logs"
LOG_FILE="$LOG_DIR/workbooksd.log"
RUNTIME_DIR="$CONSOLE_HOME/Library/Application Support/sh.workbooks.workbooksd"

# 1. Create the LaunchAgent + log directories owned by the desktop
# user (postinstall runs as root, but launchd needs the plist + log
# dir owned by the agent's UID).
sudo -u "$CONSOLE_USER" /bin/sh -c "mkdir -p '$PLIST_DIR' '$LOG_DIR'"

# Write the plist owned by the user.
sudo -u "$CONSOLE_USER" /bin/sh -c "cat > '$PLIST'" <<EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>sh.workbooks.workbooksd</string>
  <key>ProgramArguments</key>
  <array>
    <string>$BIN</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>$LOG_FILE</string>
  <key>StandardErrorPath</key><string>$LOG_FILE</string>
  <key>ProcessType</key><string>Background</string>
</dict>
</plist>
EOF

# 2. Wipe stale runtime.json BEFORE bootstrap so the new daemon's
# port-discovery write doesn't race with whatever the old one left.
rm -f "$RUNTIME_DIR/runtime.json" || true

# Bootstrap (preinstall already booted out the old one).
launchctl bootstrap "gui/$CONSOLE_UID" "$PLIST" >/dev/null 2>&1 || {
  echo "[postinstall] launchctl bootstrap failed; user may need to log out/in" >&2
}

# 3. Launch Services refresh so the bundle's CFBundleDocumentTypes
# + UTExportedTypeDeclarations are registered immediately.
LSREGISTER="/System/Library/Frameworks/CoreServices.framework/Versions/A/Frameworks/LaunchServices.framework/Versions/A/Support/lsregister"
[ -x "$LSREGISTER" ] && "$LSREGISTER" -f /Applications/Workbooks.app >/dev/null 2>&1 || true

# 4. Set Workbooks as the user's default handler for sh.workbooks.workbook.
# This writes ~/Library/Preferences/com.apple.LaunchServices/com.apple.launchservices.secure.plist
# via cfprefsd (defaults command) so the change takes immediate effect.
# We dedupe — if an entry for sh.workbooks.workbook already exists
# (e.g. a re-install), don't append a duplicate.
sudo -u "$CONSOLE_USER" /bin/sh <<INNER
  if ! defaults read com.apple.LaunchServices/com.apple.launchservices.secure LSHandlers 2>/dev/null \\
       | grep -q '"sh.workbooks.workbook"'; then
    defaults write com.apple.LaunchServices/com.apple.launchservices.secure LSHandlers \\
      -array-add '{ LSHandlerContentType = "sh.workbooks.workbook"; LSHandlerRoleEditor = "sh.workbooks.launcher"; LSHandlerPreferredVersions = { LSHandlerRoleEditor = "-"; }; }' \\
      >/dev/null 2>&1 || true
  fi
INNER

# Reload LaunchServices so the new default is picked up without a
# logout. -kill -r -domain user re-scans the user's bundle DB.
[ -x "$LSREGISTER" ] && "$LSREGISTER" -kill -r -domain user >/dev/null 2>&1 || true

exit 0
