# packages/workbooksd/release/installer/scripts/postinstall
#!/bin/sh
# Postinstall — runs as ROOT after the .pkg lays down:
#
#   /Applications/Workbooks.app   (Tauri Manager UI, opened on demand)
#   /usr/local/bin/workbooksd     (the daemon binary, supervised below)
#
# Three responsibilities:
#
#   1. Bootstrap the user-scoped LaunchAgent so the daemon runs at
#      login + restarts on crash. The user never has to open
#      Workbooks.app for the daemon to start; it's there waiting
#      whenever a browser-served workbook needs autosave or any
#      .workbook.html is double-clicked from Finder.
#
#   2. Wipe stale runtime.json from any previous (possibly
#      crashed) daemon so the freshly-bootstrapped one writes its
#      own port without stomping on another reader.
#
#   3. Refresh Launch Services so /Applications/Workbooks.app's
#      sh.workbooks.workbook UTI + public.html Alternate-rank
#      registrations land without requiring a logout.
#
# Idempotent — safe to re-run on upgrade.

set -e

CONSOLE_UID=$(stat -f%u /dev/console)
CONSOLE_USER=$(id -un "$CONSOLE_UID" 2>/dev/null || echo root)
CONSOLE_HOME=$(dscl . -read /Users/"$CONSOLE_USER" NFSHomeDirectory 2>/dev/null | awk '{print $2}')
[ -z "$CONSOLE_HOME" ] && CONSOLE_HOME="/Users/$CONSOLE_USER"

BIN="/usr/local/bin/workbooksd"
PLIST_DIR="$CONSOLE_HOME/Library/LaunchAgents"
PLIST="$PLIST_DIR/sh.workbooks.workbooksd.plist"
LOG_DIR="$CONSOLE_HOME/Library/Logs"
LOG_FILE="$LOG_DIR/workbooksd.log"
RUNTIME_DIR="$CONSOLE_HOME/Library/Application Support/sh.workbooks.workbooksd"

# Create the LaunchAgent + log directories owned by the desktop
# user (postinstall runs as root, but launchd permissions require
# the plist + log dir be readable to the agent's UID).
sudo -u "$CONSOLE_USER" /bin/sh -c "mkdir -p '$PLIST_DIR' '$LOG_DIR'"

# Write the plist owned by the user. KeepAlive=true means launchd
# restarts the daemon if it crashes; RunAtLoad=true means it starts
# immediately on bootstrap (and on every login thereafter).
sudo -u "$CONSOLE_USER" /bin/sh -c "cat > '$PLIST'" <<EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>sh.workbooks.workbooksd</string>
  <key>ProgramArguments</key>
  <array>
    <string>$BIN</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>$LOG_FILE</string>
  <key>StandardErrorPath</key><string>$LOG_FILE</string>
  <key>ProcessType</key><string>Background</string>
</dict>
</plist>
EOF

# Bootstrap (best-effort — bootout fails silently if not registered).
# Wipe stale runtime.json BEFORE bootstrap so the new daemon's port
# discovery write isn't racing with whatever the old one left.
rm -f "$RUNTIME_DIR/runtime.json" || true
launchctl bootout "gui/$CONSOLE_UID/sh.workbooks.workbooksd" >/dev/null 2>&1 || true
launchctl bootstrap "gui/$CONSOLE_UID" "$PLIST" >/dev/null 2>&1 || {
  echo "[postinstall] launchctl bootstrap failed; user may need to log out/in" >&2
}

# Launch Services refresh so /Applications/Workbooks.app's UTI
# declarations (sh.workbooks.workbook + public.html alternate) are
# picked up without a logout.
LSREGISTER="/System/Library/Frameworks/CoreServices.framework/Versions/A/Frameworks/LaunchServices.framework/Versions/A/Support/lsregister"
[ -x "$LSREGISTER" ] && "$LSREGISTER" -f /Applications/Workbooks.app >/dev/null 2>&1 || true

exit 0
